# patch/2026.05.02-157/files/webhook-router.service.example
[Unit]
Description=webhook-router
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/webhook-router
EnvironmentFile=/opt/webhook-router/.env
ExecStart=/usr/bin/python3 /opt/webhook-router/app/server.py
Restart=always
RestartSec=3
User=www-data
Group=www-data

[Install]
WantedBy=multi-user.target
